# btrfs-backup.target
[Unit]
Description=Run btrfs-backup snapshot generation and upload
Requires=btrfs-snapshot.service
Requires=btrfs-send.service
BindsTo=btrfs-send.service
