# btrfs-backup.timer
[Unit]
Description=Run btrfs-backup target

[Timer]
OnCalendar=daily
Persistent=true
Unit=btrfs-backup.target

[Install]
WantedBy=timers.target
