# btrfs-send.service
[Unit]
Requires=local-fs.target
After=local-fs.target
After=btrfs-snapshot.service

[Service]
Type=oneshot
User=emlun
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python /home/emlun/dev/btrfs-backup/btrfs-backup.py send
